# client/remotepower-agent.service
[Unit]
Description=RemotePower Agent - Remote device management client
Documentation=https://github.com/tyxak/remotepower
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/remotepower-agent run
Restart=always
RestartSec=30

# Logs to journald: journalctl -u remotepower-agent -f
StandardOutput=journal
StandardError=journal

# Security hardening
User=root
# ProtectSystem=strict  # Disabled: needs to write to /etc/remotepower/credentials
PrivateTmp=yes
# NoNewPrivileges=yes   # Disabled: breaks apt's seteuid transition to _apt (uid 105)
#                       # during exec: package-manager commands. The agent runs as
#                       # root to execute admin tasks, so this hardening is mostly
#                       # cosmetic anyway.
ProtectKernelTunables=yes
# ProtectKernelModules=yes  # NEVER re-enable: it makes /usr/lib/modules
#                           # inaccessible to the unit, so a package upgrade run
#                           # through the exec channel that triggers
#                           # update-initramfs builds an initrd with NO kernel
#                           # modules (mkinitramfs warns but exits 0) — the host
#                           # then can't find its LVM/disk root at next boot and
#                           # drops to the initramfs shell, unbootable.
ProtectControlGroups=yes

# Environment overrides (optional)
# Environment="RP_INTERVAL=60"
# v4.5.0: optional env file written by install-client.sh --ca-fingerprint; sets
# RP_CA_BUNDLE=/etc/remotepower/ca.crt so the agent trusts a self-signed CA. The
# leading '-' makes it optional (absent on real-cert / Let's Encrypt installs).
EnvironmentFile=-/etc/remotepower/agent.env

[Install]
WantedBy=multi-user.target
