# server/conf/remotepower-wsgi.service
# RemotePower persistent WSGI app tier (v5.5.0 keystone, EXPERIMENTAL — opt-in).
#
# Runs the SAME api.py under gunicorn (persistent workers) instead of the
# fork-per-request CGI / SCGI models. Validated correct under load on both the
# SQLite and Postgres backends (per-request isolation holds: no load-cache /
# correlation-id / response-body bleed). The default + fully supported deployment
# is still CGI (api_cgi.py) or the SCGI worker (remotepower-api.service); use this
# only for large fleets that want a pre-warmed persistent tier, and keep CGI as the
# fallback until you've validated it under your own load. See docs/wsgi.md.
#
# Process model: request state (context, output buffer, DB connections) is
# thread-local, so the worker serves requests CONCURRENTLY on threads — gthread is
# supported. Tune --workers (processes) x --threads to your CPU + DB-pool budget.
#
# Prereqs: pip install gunicorn  (gunicorn is NOT a RemotePower dependency).
# Install:
#   cp server/conf/remotepower-wsgi.service /etc/systemd/system/
#   systemctl daemon-reload && systemctl enable --now remotepower-wsgi
# then point nginx /api/ at http://127.0.0.1:8090 (proxy_pass; see docs/wsgi.md)
# and reload nginx. Pair with remotepower-scheduler.service + RP_EXTERNAL_SCHEDULER=1
# (in api.env) so the cadence runs out-of-band (measured ~25x lower request latency
# on a networked DB). Roll back by repointing nginx at fcgiwrap/SCGI.

[Unit]
Description=RemotePower WSGI app tier (gunicorn, experimental)
Documentation=file:///var/www/remotepower/cgi-bin/wsgi.py
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/remotepower/cgi-bin
# NB: deliberately NO `RuntimeDirectory=remotepower`. gunicorn binds TCP
# (127.0.0.1:8090) and needs nothing in /run. If this unit declared the shared
# `remotepower` RuntimeDirectory, systemd would DELETE /run/remotepower when this
# unit stops/restarts — which would nuke the SCGI worker's /run/remotepower/api.sock
# out from under a co-resident remotepower-api during a CGI⇄WSGI cutover.
Environment=RP_DATA_DIR=/var/lib/remotepower
EnvironmentFile=-/etc/remotepower/api.env
# Threaded workers (request state is thread-local). Bind to localhost for nginx.
ExecStart=/usr/bin/gunicorn --workers 4 --threads 8 --timeout 120 \
          --bind 127.0.0.1:8090 wsgi:application
Restart=always
RestartSec=2

# Hardening — mirrors remotepower-api.service (same cadence/feature surface).
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/remotepower
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
Environment=RP_WG_DIRECT=1
Environment=RP_WG_DIR=/var/lib/remotepower/wg

[Install]
WantedBy=multi-user.target
